# contrib/systemd/ryzen-stabilizator.target
[Unit]
Description=Ryzen Stabilizator Tabajara
Requires=ryzen-stabilizator.service ryzen-stabilizator@boot.service ryzen-stabilizator@resume.service
